# plocate-updatedb.service.in
[Unit]
Description=Update the plocate database
ConditionACPower=true
RandomizedDelaySec=12h
AccuracySec=20min

[Service]
Type=oneshot
ExecStart=@sbindir@/@updatedb_progname@
LimitNOFILE=131072
IOSchedulingClass=idle

PrivateTmp=true
PrivateDevices=true
PrivateNetwork=true
