# plocate-updatedb.service.in
[Unit]
Description=Update the plocate database
ConditionACPower=true

[Service]
Type=oneshot
ExecStart=@sbindir@/@updatedb_progname@
LimitNOFILE=131072
IOSchedulingClass=idle
Nice=19

PrivateTmp=true
PrivateDevices=true
PrivateNetwork=true
